# systemd/skybridge.service
[Unit]
Description=PX4 SkyBridge Relay
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/groove/work/PX4-SkyBridge-Relay
ExecStart=/usr/bin/env bash /home/groove/work/PX4-SkyBridge-Relay/scripts/bridge_supervisor.sh
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
